# mautrix-signal/mautrix-signal.service
[Unit]
Description=Matrix-signal bridge
After=synapse.service

[Service]
Restart=always

ExecStart=/usr/bin/python3 -m mautrix_signal -c /etc/mautrix/signal/config.yaml -r /etc/mautrix/signal/registration.yaml

WorkingDirectory=/var/lib/mautrix/signal

StateDirectory=mautrix/signal

DynamicUser=true

PrivateTmp=yes
PrivateDevices=true
PrivateUsers=true

DevicePolicy=closed

ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectClock=true
ProtectKernelLogs=yes

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

NoNewPrivileges=yes
LockPersonality=yes

SystemCallFilter=~@mount
[Install]
WantedBy=multi-user.target
