# mautrix-instagram/mautrix-instagram.service
[Unit]
Description=Matrix-instagram bridge
After=synapse.service

[Service]
Restart=always

LoadCredential=config:/etc/mautrix/instagram/config.yaml
ExecStart=/usr/bin/python3 -m mautrix_instagram -c ${CREDENTIALS_DIRECTORY}/config

WorkingDirectory=/var/lib/mautrix/instagram

StateDirectory=mautrix/instagram

DynamicUser=true

PrivateTmp=yes
PrivateDevices=true
PrivateUsers=true

DevicePolicy=closed

ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectClock=true
ProtectKernelLogs=yes

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

NoNewPrivileges=yes
LockPersonality=yes

SystemCallFilter=~@mount
[Install]
WantedBy=multi-user.target
